# unzipper.timer
[Unit]
Description=Run Automated Unzipper for Downloads every 15 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min
Unit=unzipper.service

[Install]
WantedBy=timers.target
